# boiler-bot.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/boiler/boiler-bot.sh

[Install]
WantedBy=default.target
